# src/crontab
0 8 * * * cd /app; /usr/local/bin/python3 /app/daily.py
